# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn --bind :8080 --workers 3 run:app
directory=/app
autostart=true
autorestart=true
redirect_stderr=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
redirect_stderr=true
